# infrastructure/gpu-worker.service
[Unit]
Description=GPU Worker Container Wrapper
# Wait for Docker and Network to be ready before running
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=simple
# Point to the script we just made
ExecStart=/opt/run-worker.sh
# Do NOT restart automatically; the script handles the logic
Restart=no
User=root

[Install]
WantedBy=multi-user.target